# Firmware/system/br2_external/board/brass/linux/rootfs-overlay/etc/profile.d/python.sh
export PYTHONPATH="/usr/share/brass/python"
